# .config/corsair/corsair-headset-fix.service
[Unit]
Description=Fix Corsair VOID ELITE headset audio (analog profile + gain)
After=pipewire.service wireplumber.service
Wants=pipewire.service wireplumber.service

[Service]
Type=oneshot
ExecStart=%h/.config/corsair/corsair-headset-fix

[Install]
WantedBy=default.target
